# system/skeleton/etc/profile.d/cmsmaddrs.sh
read -r CM < /etc/cm_addr
read -r SM < /etc/sm_addr

export CM
export SM
